# docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:api]
command=/opt/venv/bin/python3 -m uvicorn api.main:app --host 127.0.0.1 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ui]
command=/opt/venv/bin/python3 -m streamlit run ui/app.py --server.port 8501 --server.address 127.0.0.1 --server.headless true
directory=/app
autostart=true
autorestart=true
environment=API_URL="http://127.0.0.1:8000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
